# etc/profile.d/welcome.sh
deersay "Welcome to DeerPi for Orange pi zero"
